# contrib/systemd/docker-volume-glusterfs.socket
[Unit]
Description=Gluster volume driver for Docker

[Socket]
ListenStream=/run/docker/plugins/glusterfs.sock

[Install]
WantedBy=sockets.target
